# resources/rt-app.service
[Unit]
Description=Request Tracker Search Application
After=network.target elasticsearch.service
Wants=elasticsearch.service

[Service]
ExecStart=/usr/bin/node /home/ma/repo/rt-app/resources/html/server/rt-app.js
WorkingDirectory=/home/ma/repo/server
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=rt-app
User=ma
Group=ma
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
